# systemd/BookWL.service
[Unit]
Description=BookWL Service
Wants=BookWL.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/free

[Install]
WantedBy=multi-user.target